# .config/systemd/user/autocutsel.service
[Unit]
Description=autocutsel
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/autocutsel

[Install]
WantedBy=graphical-session.target

