# .config/systemd/user/polybar@.service
[Unit]
Description=Polybar
Requires=graphical-session.target
After=graphical-session.target pulseaudio.socket i3.service

[Service]
Type=simple
ExecStart=/usr/bin/polybar --reload --log=warning %i
Restart=always

[Install]
WantedBy=desktop.target
