# systemd/user/blueman-applet.service
[Unit]
Description=Bluetooth Applet
After=graphical-session.target
Requires=graphical-session.target
ConditionPathExists=/sys/class/bluetooth

[Service]
ExecStart=/usr/bin/blueman-applet
Restart=on-failure
Slice=background.slice

[Install]
WantedBy=desktop.target
