# systemd/user/xorg.target
[Unit]
Description=Graphical desktop environment
Requires=basic.target
