# systemd/user/hyprpaper.service.d/override.conf
[Service]
Slice=session-graphical.slice
