# systemd/user/dunst.service.d/override.conf
[Service]
Slice=session-graphical.slice

[Install]
WantedBy=graphical-session.target
